# install/proxpilot.service
[Unit]
Description=ProxPilot — GitOps VM lifecycle manager for Proxmox
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/proxpilot -config /opt/proxpilot/config.yml
Restart=always
RestartSec=10
WorkingDirectory=/opt/proxpilot

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proxpilot

[Install]
WantedBy=multi-user.target
